# infrastructure/roles/tekton/files/k0s.service
[Unit]
Description=Start k0s

[Service]
Type=simple
ExecStart=/usr/bin/k0s controller
Restart=on-failure
